# build-armbian/armbian-files/different-files/ai3399-c/rootfs/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

bash /etc/custom_service/start_service.sh
# ai3399-c watchdog start. dearxjoe
chmod +x /usr/bin/mcu-watchdog.sh
nohup /usr/bin/mcu-watchdog.sh &
# ai3399-c watchdog   end.
exit 0
